# systemd/bbalert-prod.service
[Unit]
Description=BBAlert Telegram Bot - Production Environment
Documentation=https://github.com/ersus93/bbalert
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=USUARIO
Group=USUARIO
WorkingDirectory=/home/USUARIO/bbalert-prod/bbalert

# Comando de ejecución
ExecStart=/home/USUARIO/bbalert-prod/venv/bin/python bbalert.py

# Variables de entorno
Environment="PYTHONUNBUFFERED=1"
Environment="PATH=/home/USUARIO/bbalert-prod/venv/bin:%h/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Reinicio automático
Restart=on-failure
RestartSec=5s

# Configuración de logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bbalert-prod

# Límites de recursos (Según docs/limpieza.md)
MemoryMax=150M
MemorySwapMax=50M
LimitNOFILE=1024
TimeoutStartSec=30
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
